# data/supervisor/stoqserver.conf
[program:stoqserver]
directory=/var/lib/stoqserver
command=stoqserver run
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/stoqserver.log
stdout_logfile_backups=1000
environment=PGPASSFILE=/var/lib/stoqserver/.pgpass
autostart=true
autorestart=true
startretries=100
